# deploy/systemd/pki-ra-engine.service
[Unit]
Description=PKI RA Engine
Documentation=https://github.com/straptrust/pki
After=network.target postgresql.service pki-ca-engine.service
Requires=postgresql.service
Wants=pki-ca-engine.service

[Service]
Type=exec
User=pki
Group=pki
WorkingDirectory=/opt/pki/releases/ra_engine

EnvironmentFile=/opt/pki/.env
Environment=PHX_SERVER=true
Environment=RELEASE_DISTRIBUTION=name
Environment=RELEASE_NODE=pki_ra_engine@127.0.0.1
Environment=RELEASE_COOKIE_FILE=/opt/pki/.cookies/ra_engine
Environment=ELIXIR_ERL_OPTIONS=+S 2:2 +SDcpu 2

ExecStart=/opt/pki/releases/ra_engine/bin/pki_ra_engine start
ExecStop=/opt/pki/releases/ra_engine/bin/pki_ra_engine stop

Restart=on-failure
RestartSec=5
TimeoutStartSec=60
TimeoutStopSec=30

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/pki/releases/ra_engine /tmp

StandardOutput=journal
StandardError=journal
SyslogIdentifier=pki-ra-engine

[Install]
WantedBy=multi-user.target
