# deploy/systemd/pki-engines.service
[Unit]
Description=PKI Engines (Audit Trail receiver + PlatformEngine bootstrap modules)
Documentation=https://github.com/straptrust/pki
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=exec
User=pki
Group=pki
WorkingDirectory=/opt/pki/releases/engines

EnvironmentFile=/opt/pki/.env
# CA/RA engines no longer expose HTTP APIs, and Validation (OCSP/CRL/TSA)
# runs inside each tenant BEAM (via pki_tenant_node). This release only
# hosts the centralized audit-trail receiver + platform bootstrap modules.
Environment=RELEASE_DISTRIBUTION=name
Environment=RELEASE_NODE=pki_engines@127.0.0.1
Environment="ELIXIR_ERL_OPTIONS=+S 2:2 +SDcpu 2:2"

ExecStart=/opt/pki/releases/engines/bin/pki_engines start
ExecStop=/opt/pki/releases/engines/bin/pki_engines stop

Restart=on-failure
RestartSec=10
RestartMaxDelaySec=120
StartLimitIntervalSec=300
StartLimitBurst=5
TimeoutStartSec=90
TimeoutStopSec=30

# Basic hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/pki/releases/engines /var/lib/softhsm /tmp

StandardOutput=journal
StandardError=journal
SyslogIdentifier=pki-engines

[Install]
WantedBy=multi-user.target
