# deploy/systemd/pki-engines.service
[Unit]
Description=PKI Engines (audit-trail receiver + platform bootstrap modules)
Documentation=https://github.com/straptrust/pki
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=exec
User=pki
Group=pki
WorkingDirectory=/opt/pki/releases/engines

EnvironmentFile=/opt/pki/.env
# CA/RA engines no longer expose HTTP APIs, and validation (OCSP/CRL)
# runs inside each tenant BEAM. This release hosts the centralized
# audit-trail receiver and platform bootstrap modules only.
Environment=RELEASE_DISTRIBUTION=name
Environment=RELEASE_NODE=pki_engines@127.0.0.1

# ── Mnesia ────────────────────────────────────────────────────────────────────
Environment=MNESIA_DIR=/var/lib/pki/mnesia/engines

# ── BEAM runtime flags ────────────────────────────────────────────────────────
# +S 0:0      — use all cores; engines is lightweight but don't waste headroom
# +SDcpu 1:1  — only 1 dirty-CPU scheduler; minimal NIF usage in this release
# +SDio 4:4   — 4 dirty-I/O schedulers for Mnesia and audit log writes
# +A 16       — 16-thread async pool; smaller than platform, adequate for audit
# +sbwt none +swt very_low — don't spin schedulers; engines is mostly idle
# +MBas aobf +MBlmbcs 512  — allocator tuning (same as platform)
Environment="ELIXIR_ERL_OPTIONS=+S 0:0 +SDcpu 1:1 +SDio 4:4 +A 16 +sbwt none +swt very_low +MBas aobf +MBlmbcs 512"

Environment=ERL_FULLSWEEP_AFTER=20

ExecStart=/opt/pki/releases/engines/bin/pki_engines start
ExecStop=/opt/pki/releases/engines/bin/pki_engines stop

Restart=on-failure
RestartSec=10
RestartMaxDelaySec=120
StartLimitIntervalSec=300
StartLimitBurst=5
TimeoutStartSec=90
TimeoutStopSec=30

# ── Resource limits ───────────────────────────────────────────────────────────
LimitNOFILE=65536
LimitNPROC=65536
OOMScoreAdjust=-400

# ── Hardening ─────────────────────────────────────────────────────────────────
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/pki/releases/engines /var/lib/pki/mnesia/engines /var/lib/softhsm /tmp

StandardOutput=journal
StandardError=journal
SyslogIdentifier=pki-engines

[Install]
WantedBy=multi-user.target
