# deploy/systemd/pki-audit.service
[Unit]
Description=PKI Audit Trail Service
Documentation=https://github.com/straptrust/pki
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=exec
User=pki
Group=pki
WorkingDirectory=/opt/pki/releases/audit

EnvironmentFile=/opt/pki/.env
Environment=RELEASE_DISTRIBUTION=name
Environment=RELEASE_NODE=pki_audit@127.0.0.1

# ── Mnesia ────────────────────────────────────────────────────────────────────
Environment=MNESIA_DIR=/var/lib/pki/mnesia/audit

# ── BEAM runtime flags ────────────────────────────────────────────────────────
# Audit is deliberately lightweight — 1 scheduler core, minimal async pool.
# +sbwt none prevents wasted CPU on a process that handles only sporadic writes.
Environment="ELIXIR_ERL_OPTIONS=+S 1:1 +SDcpu 1:1 +A 8 +sbwt none +swt very_low"

Environment=ERL_FULLSWEEP_AFTER=20

ExecStart=/opt/pki/releases/audit/bin/pki_audit start
ExecStop=/opt/pki/releases/audit/bin/pki_audit stop

Restart=on-failure
RestartSec=10
RestartMaxDelaySec=120
StartLimitIntervalSec=300
StartLimitBurst=5
TimeoutStartSec=90
TimeoutStopSec=30

# ── Resource limits ───────────────────────────────────────────────────────────
LimitNOFILE=32768
LimitNPROC=32768
OOMScoreAdjust=-300

# ── Hardening ─────────────────────────────────────────────────────────────────
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/pki/releases/audit /var/lib/pki/mnesia/audit /tmp

StandardOutput=journal
StandardError=journal
SyslogIdentifier=pki-audit

[Install]
WantedBy=multi-user.target
