# deploy/systemd/pki-ra-portal.service
[Unit]
Description=PKI RA Portal
Documentation=https://github.com/straptrust/pki
After=network.target pki-ra-engine.service
Wants=pki-ra-engine.service

[Service]
Type=exec
User=pki
Group=pki
WorkingDirectory=/opt/pki/releases/ra_portal

EnvironmentFile=/opt/pki/.env
Environment=PHX_SERVER=true
Environment=RELEASE_DISTRIBUTION=name
Environment=RELEASE_NODE=pki_ra_portal@127.0.0.1
Environment=RELEASE_COOKIE_FILE=/opt/pki/.cookies/ra_portal
Environment=ELIXIR_ERL_OPTIONS=+S 2:2 +SDcpu 2

ExecStart=/opt/pki/releases/ra_portal/bin/pki_ra_portal start
ExecStop=/opt/pki/releases/ra_portal/bin/pki_ra_portal stop

Restart=on-failure
RestartSec=5
TimeoutStartSec=60
TimeoutStopSec=30

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/pki/releases/ra_portal /tmp

StandardOutput=journal
StandardError=journal
SyslogIdentifier=pki-ra-portal

[Install]
WantedBy=multi-user.target
