# deploy/systemd/pki-ca-portal.service
[Unit]
Description=PKI CA Portal
Documentation=https://github.com/straptrust/pki
After=network.target pki-ca-engine.service
Wants=pki-ca-engine.service

[Service]
Type=exec
User=pki
Group=pki
WorkingDirectory=/opt/pki/releases/ca_portal

EnvironmentFile=/opt/pki/.env
Environment=PHX_SERVER=true
Environment=RELEASE_DISTRIBUTION=name
Environment=RELEASE_NODE=pki_ca_portal@127.0.0.1
Environment="ELIXIR_ERL_OPTIONS=+S 2:2 +SDcpu 2:2"

ExecStart=/opt/pki/releases/ca_portal/bin/pki_ca_portal start
ExecStop=/opt/pki/releases/ca_portal/bin/pki_ca_portal stop

Restart=on-failure
RestartSec=5
TimeoutStartSec=60
TimeoutStopSec=30

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/pki/releases/ca_portal /tmp

StandardOutput=journal
StandardError=journal
SyslogIdentifier=pki-ca-portal

[Install]
WantedBy=multi-user.target
