# deploy/systemd/pki-portals.service
[Unit]
Description=PKI Portals (CA Portal, RA Portal, Platform Portal + engines in direct mode)
Documentation=https://github.com/straptrust/pki
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=exec
User=pki
Group=pki
WorkingDirectory=/opt/pki/releases/portals

EnvironmentFile=/opt/pki/.env
Environment=PHX_SERVER=true
Environment=ENGINE_CLIENT_MODE=direct
Environment=RELEASE_DISTRIBUTION=name
Environment=RELEASE_NODE=pki_portals@127.0.0.1
# Portals + in-process engines — moderate scheduler count
Environment="ELIXIR_ERL_OPTIONS=+S 2:2 +SDcpu 2:2"

ExecStart=/opt/pki/releases/portals/bin/pki_portals start
ExecStop=/opt/pki/releases/portals/bin/pki_portals stop

Restart=on-failure
RestartSec=5
TimeoutStartSec=60
TimeoutStopSec=30

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/pki/releases/portals /var/lib/softhsm /tmp

StandardOutput=journal
StandardError=journal
SyslogIdentifier=pki-portals

[Install]
WantedBy=multi-user.target
