# deploy/systemd/pki-platform-portal.service
[Unit]
Description=PKI Platform Admin Portal
Documentation=https://github.com/straptrust/pki
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=exec
User=pki
Group=pki
WorkingDirectory=/opt/pki/releases/platform_portal

EnvironmentFile=/opt/pki/.env
Environment=PHX_SERVER=true
Environment=RELEASE_DISTRIBUTION=name
Environment=RELEASE_NODE=pki_platform_portal@127.0.0.1
Environment=RELEASE_COOKIE_FILE=/opt/pki/.cookies/platform_portal
Environment=ELIXIR_ERL_OPTIONS=+S 2:2 +SDcpu 2

ExecStart=/opt/pki/releases/platform_portal/bin/pki_platform_portal start
ExecStop=/opt/pki/releases/platform_portal/bin/pki_platform_portal stop

Restart=on-failure
RestartSec=5
TimeoutStartSec=60
TimeoutStopSec=30

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/pki/releases/platform_portal /tmp

StandardOutput=journal
StandardError=journal
SyslogIdentifier=pki-platform-portal

[Install]
WantedBy=multi-user.target
