# deploy/systemd/pki-audit.service
[Unit]
Description=PKI Audit Trail Service
Documentation=https://github.com/straptrust/pki
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=exec
User=pki
Group=pki
WorkingDirectory=/opt/pki/releases/audit

EnvironmentFile=/opt/pki/.env
Environment=RELEASE_DISTRIBUTION=name
Environment=RELEASE_NODE=pki_audit@127.0.0.1
# Lightweight service — minimal schedulers
Environment="ELIXIR_ERL_OPTIONS=+S 1:1 +SDcpu 1:1"

ExecStart=/opt/pki/releases/audit/bin/pki_audit start
ExecStop=/opt/pki/releases/audit/bin/pki_audit stop

Restart=on-failure
RestartSec=5
TimeoutStartSec=60
TimeoutStopSec=30

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/pki/releases/audit /tmp

StandardOutput=journal
StandardError=journal
SyslogIdentifier=pki-audit

[Install]
WantedBy=multi-user.target
